# Debian_Server_Init/new_version/other/multi_nic_route/multi-nic-route.service
[Unit]
Description=Multi-NIC Route Configuration Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=【path to】/multi_nic_route.sh
Restart=no

[Install]
WantedBy=multi-user.target

