# Debian_GNOME_Init/other/multi_nic_route/multi-nic-route.service
[Unit]
Description=Multi-NIC Route Configuration Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# 请将下方路径改为脚本实际所在目录的绝对路径，或先将脚本复制到 /usr/local/bin 后使用 /usr/local/bin/multi_nic_route.sh
ExecStart=/usr/local/bin/multi_nic_route.sh
Restart=no

[Install]
WantedBy=multi-user.target

